# packaging/systemd/kaizen-daemon.service
[Unit]
Description=Kaizen local daemon
After=network.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/kaizen daemon start
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
